# roles/photogram/files/upstart/photogram-api-pictures.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

setuid root
chdir /opt/photogram-api

env NODE_ENV="production"
env PHOTOGRAM_SECRET="b198eded-3e56-4db9-a711-cddd807585b7"

exec npm run pictures
